# src/prefect/agent/local/supervisord.conf
; Supervisor config file for running Prefect Agents.
;
; To run for the first time:
; Save this file with the name `supervisord.conf` and run:
; supervisord # will start all programs
; supervisorctl # will open an interactive prompt for inspection
; 
; If you need to edit this file and restart:
; supervisorctl reload # will reload config and restart programs
;
; For more information on configuring supervisor, please see:
; http://supervisord.org/configuration.html
;
; To configure supervisor logging, please see:
; http://supervisord.org/logging.html
;
; To inspect the status of a given program, 
; use the interactive supervisorctl command:
; supervisorctl
; fg prefect-agent
;
; To stop all running programs, run:
; supervisorctl stop all
;
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
loglevel=debug               ; log level; default info; others: debug,warn,trace

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:prefect-agent]
command=prefect agent local start {{OPTS}}

; you can add as many programs as you'd like
; including multiple Prefect Agents
; [program:foo-agent]
; command=prefect agent start local -l foo
